# scripts/signal-cli.service
[Unit]
Description=signal-cli daemon (JSON-RPC socket mode)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Edit this: set your registered phone number
ExecStart=%h/.local/bin/signal-cli -a +44XXXXXXXXXX daemon --socket --receive-mode on-connection
Restart=on-failure
RestartSec=10

# Socket will be created at $XDG_RUNTIME_DIR/signal-cli/socket
RuntimeDirectory=signal-cli

[Install]
WantedBy=default.target
